# Rasp_PI/mqtt_publisher.service
[Unit]
Description=BLE Medicine Tag Scanner and MQTT Publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/iot_project
ExecStart=/home/pi/iot_project/venv/bin/python mqtt_publisher.py
Restart=on-failure
RestartSec=5

# Log to systemd journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
